# django-api/supervisord.conf
[unix_http_server]
file = /var/run/supervisor.sock   ; the path to the socket file


[supervisord]
nodaemon=true                 ; start in foreground if true; default false
;logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
;logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
;logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
;loglevel=info                ; log level; default info; others: debug,warn,trace
;pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
;silent=true                 ; no logs to stdout if true; default false
;minfds=1024                  ; min. avail startup file descriptors; default 1024
;minprocs=200                 ; min. avail process descriptors;default 200


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisorctl]
serverurl=unix:///var/run/supervisor.sock


[program:gunicorn_wsgi]
command=gunicorn --chdir project --bind :8001 project.wsgi:application
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile= /opt/services/djangoapp/src/log/gunicorn_supervisor_stdout.log


;[program:daphne_asgi]
;command=daphne -b 0.0.0.0 -p 8001 project.asgi:application
;directory=/opt/services/djangoapp/src/project/
;;autostart=true
;autorestart=true
;stopasgroup=true
;user=root
;stdout_logfile = /opt/services/djangoapp/src/log/daphne_supervisor_error.log



;[fcgi-program:asgi]
; TCP socket used by Nginx backend upstream
;socket=tcp://localhost:8001

; Directory where your site's project files are located
directory=/opt/services/djangoapp/src/project/

; Each process needs to have a separate socket file, so we use process_num
; Make sure to update "mysite.asgi" to match your project name
;command=daphne -u /run/daphne/daphne%(process_num)d.sock --fd 0 --access-log - --proxy-headers project.asgi:application

; Number of processes to startup, roughly the number of CPUs you have
;numprocs=1

; Give each process a unique name so they can be told apart
;process_name=asgi%(process_num)d

; Automatically start and recover processes
;autostart=true
;autorestart=true

; Choose where you want your log to go
;stdout_logfile=/opt/services/djangoapp/src/log/daphne_supervisor_stdout.log
;stderr_logfile=/opt/services/djangoapp/src/log/daphne_supervisor_error.log
;redirect_stderr=true